# etc/update-motd.d/00-header
#!/bin/sh
printf  '\e[1;33m%-6s\e[m\n' "$(hostname) (encendido desde $(uptime -s))"
printf "($(lsb_release -sd) $(lsb_release -sc) - kernel $(uname -r) $(uname -m))\n"
